# supervisord.conf
[supervisord]
nodaemon=true

[program:freshclam]
command=freshclam -d

[program:clamav]
command=clamd

[program:clamav-rest]
command=/go-clamav-rest